# xhs-ocr.service
[Unit]
Description=XHS OCR Service
After=network.target

[Service]
Type=simple
User=callfans
WorkingDirectory=/home/callfans/data/AI-box/XHS-OCR
Environment=TZ=Asia/Shanghai
ExecStart=/home/callfans/miniconda3/bin/conda run -n py312 python /home/callfans/data/AI-box/XHS-OCR/xhs_ocr.py --mode schedule --at-time "18:00"
Restart=always
StandardOutput=append:/home/callfans/data/AI-box/XHS-OCR/logs/run.log
StandardError=append:/home/callfans/data/AI-box/XHS-OCR/logs/run.log

[Install]
WantedBy=multi-user.target